# docker/etc/s6/gopls/run
#!/bin/bash

# 启动Gopls
exec su-exec "${USERNAME}" gopls -listen=":${PORT}" -logfile=auto -debug=:0 -rpc.trace serve
